# AdGuardHome/s6-overlay/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root AdGuardHome -c /AdGuardHome/AdGuardHome.yaml &>/dev/null
